# assets/duckdns.service
[Unit]
Description=Refresh DuckDNS ip address

[Service]
Type=oneshot
EnvironmentFile=${HOME}/.config/devops.env
ExecStart=${HOME}/bin/duckdns.sh
