# assets/ngrok-start.service
[Unit]
Description=Start ngrok tunnel
After=network.target

[Service]
Type=simple
User=${USER}
EnvironmentFile=${HOME}/.config/devops.env
ExecStart=${HOME}/bin/ngrok-start.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
